# conf.d/FeelingSurfViewer.conf
[program:FeelingSurfViewer]
command=timeout 30m xvfb-run FeelingSurfViewer --no-sandbox --cap-add SYS_ADMIN
autorestart=true
