# conf.d/FeelingSurfViewer.conf
[program:FeelingSurfViewer]
command=timeout 20m FeelingSurfViewer --no-sandbox --cap-add SYS_ADMIN
autorestart=true
